# rootfs/debian/features/firstboot-info/overlay/etc/update-motd.d/50-rk3588
#!/bin/sh
[ -r /var/lib/rk3588-board-info ] || exit 0
echo
echo "RK3588 board info:"
sed 's/^/  /' /var/lib/rk3588-board-info
if command -v nmtui >/dev/null 2>&1; then
    echo "  network: use nmtui or nmcli"
fi
echo
