# rootfs/debian/overlays/firstboot/overlay/etc/systemd/system/sbc-firstboot.service
[Unit]
Description=Expand RK3588 root filesystem on first boot
After=local-fs.target
ConditionPathExists=!/var/lib/sbc-firstboot.done

[Service]
Type=oneshot
ExecStart=-/usr/local/sbin/sbc-firstboot
RemainAfterExit=yes
TimeoutStartSec=10min

[Install]
WantedBy=multi-user.target
